# etc/profile.d/maven.sh
export MAVEN_HOME=/opt/apache-maven-3.2.3
export MAVEN_OPTS="-Xmx2g -XX:MaxPermSize=512M -XX:ReservedCodeCacheSize=512m -Xms256m -Xss8m"
export PATH=$PATH:$MAVEN_HOME/bin
