# etc/profile.d/spark.sh
export SPARK_HOME=/usr/local/spark-1.4.1-bin-2.2.0
